# byteshelf.service
[Unit]
Description=ByteShelf File Storage API
After=network.target

[Service]
Type=simple
User=byteshelf
WorkingDirectory=/opt/byteshelf
ExecStart=/opt/byteshelf/ByteShelf
Restart=always
RestartSec=2
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000
Environment=StoragePath=/var/byteshelf/storage
Environment=ChunkConfiguration__ChunkSizeBytes=1048576

[Install]
WantedBy=multi-user.target 